# dist/DEBIAN/postinst
#!/bin/sh
systemctl enable iot-wifi-setup.service
systemctl start iot-wifi-setup.service