# dist/DEBIAN/prerm
#!/bin/sh
systemctl stop iot-wifi-setup.service
systemctl disable iot-wifi-setup.service